# playbooks/templates/todo.service.j2
[Unit]
Description=todo

[Service]
Type=simple
Restart=always
RestartSec=5s
ExecStart={{ backend_tmp_bin }}
EnvironmentFile={{ backend_env_file }}
[Install]
WantedBy=multi-user.target